# systemd/user/battery.service
[Unit]
Description=Battery Monitor

[Service]
Type=oneshot
ExecStart=/home/shyweeds/.local/bin/battery_checker.sh
